# claw-hive.service
[Unit]
Description=CLAW-HIVE Dashboard
After=network.target

[Service]
Type=simple
User=arthur
WorkingDirectory=/home/arthur/claw-hive
ExecStart=/home/arthur/.nvm/versions/node/v22.22.0/bin/node src/server.js
Environment=PORT=3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
